# config/systemd/user/myntpd.service
[Unit]
Description= NTPD Script

[Service]
ExecStart=/usr/bin/ntpd-script

[Install]
WantedBy=default.target
